# configs/logrotate/crontab
0 3 * * * /usr/sbin/logrotate /etc/logrotate.d/nginx
